# LaunchAgents/com.zkoppert.triage-dependabot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zkoppert.triage-dependabot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/zkoppert/repos/dotfiles/bin/triage-dependabot</string>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <!-- Every hour, on the hour, 24/7. Catches dependabot bumps that land
         outside business hours (evenings, weekends) so they don't sit in
         the notifications inbox until Monday morning. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/zkoppert/Library/Logs/triage-dependabot.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/zkoppert/Library/Logs/triage-dependabot.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/zkoppert</string>
    </dict>
</dict>
</plist>
